# startup/files/systemd/user/autorandr.service
[Unit]
Description=Enact Auto configures Xrandr

[Service]
ExecStart=autorandr --watch