# startup/files/systemd/user/synergy.service
[Unit]
Description=Synergy Mouse and Keyboard Sharing

[Service]
ExecStart=synergy
Environment="DISPLAY=:1"