# windowmanagers/files/enact/systemd/user/enact.service
[Unit]
Description=Enact Auto configures Xrandr

[Service]
ExecStart=$HOME/.local/share/bin/enact --pos top --watch